# alix-templates/alix-base-1004/files/etc/update-motd.d/10-help-text
#!/bin/sh

. /etc/uictpl.conf

echo
echo "Welcome to Ubuntu!"
echo
echo "Installed from Ubuntu (${UIC_RELEASE}) on ${UIC_CREATION_TIMESTAMP} using"
echo " - UIC (Ubuntu Installation Creator) ${UIC_CREATION_VERSION}"
echo " - ${UIC_SRCNAME} ${UIC_SRCVERSION} uic template"
echo
